# services/mysqld/run
#!/bin/execlineb -P
s6-envuidgid service
if { s6-mkdir -p /run/mysqld }
if { s6-chown -U /run/mysqld }
s6-applyuidgid -U -z
mysqld
