# services/transmission/run
#!/bin/execlineb -P
s6-envuidgid transmission
s6-applyuidgid -U -z
transmission-daemon -f -g /var/lib/transmission
