# services/.s6-svscan/finish
#!/bin/execlineb -P
if { s6-nuke -t }
wait -t 5000 { }
