# services/smtpd/run
#!/bin/execlineb -P
s6-envuidgid smtpd
s6-applyuidgid -U -z
smtpd -d
