# services/smtpd/run
#!/bin/execlineb -P
s6-envuidgid service
s6-applyuidgid -U -z
smtpd -d
