# services/sshd/run
#!/bin/execlineb -P
if { ssh-keygen -A }
/bin/sshd -De
