# services/murmur/run
#!/bin/execlineb -P
s6-envuidgid murmur
s6-applyuidgid -U -z
murmur -fg -ini /etc/murmur/murmur.ini
